# borgbackup/server/drive-insertion/etc/backups/automatic-backup.service
[Service]
Type=oneshot
ExecStart=/etc/backups/run.sh
